# systemd/outage-monitor.service
[Unit]
Description=Internet connectivity outage monitor (DNS + LAN ping to MySQL)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust paths and user for your deployment.
# User=outage
# Group=outage
WorkingDirectory=/home/kushagra/Documents/projects/internet-outage
EnvironmentFile=-/etc/outage-monitor.env
ExecStart=/home/kushagra/Documents/projects/internet-outage/.venv/bin/python -m outage_monitor > /home/kushagra/Documents/logs/outage_monitor.log
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
